# gov-fetch.service
[Unit]
Description=Gov Fetch Service
Documentation=https://github.com/arunabhcode/gov-fetch
After=docker.service network-online.target
Requires=docker.service network-online.target

[Service]
Type=simple
User=%u
Group=%g
# If GOVFETCH_DIR is not set, use a default
Environment="GOVFETCH_DIR=/opt/gov-fetch"
# Auto-detect profile based on GPU presence (will be overridden by ExecStart)
Environment="COMPOSE_PROFILES=cpu"
WorkingDirectory=/opt/gov-fetch

# Auto-detect if GPU is present and set COMPOSE_PROFILES accordingly
ExecStartPre=/bin/bash -c "${GOVFETCH_DIR}/scripts/detect-gpu.sh"
ExecStartPre=/usr/bin/docker compose -f ${GOVFETCH_DIR}/docker-compose.yml down
ExecStart=/usr/bin/docker compose -f ${GOVFETCH_DIR}/docker-compose.yml --profile ${COMPOSE_PROFILES} up
ExecStop=/usr/bin/docker compose -f ${GOVFETCH_DIR}/docker-compose.yml down

Restart=always
RestartSec=10s
TimeoutStartSec=180
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target